# config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:lstn]
command=/usr/local/bin/uwsgi --plugin python /etc/uwsgi.ini
stdout_logfile=/var/log/lstn.log
redirect_stderr=true
